# ferm.service.in
[Unit]
Description=for Easy Rule Making

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@prefix@/sbin/ferm /etc/ferm.conf
ExecStop=@prefix@/sbin/ferm -F /etc/ferm.conf

[Install]
WantedBy=multi-user.target
